# infrastructure/docker-logs-s3/systemd/db-logs-export.service
[Unit]
Description=Export Docker database logs to S3
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
ExecStart=/opt/docker-logs-s3/scripts/export-db-logs.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=db-logs-export

[Install]
WantedBy=multi-user.target
